# video_recorder.service
[Unit]
Description=Wedding Phone Video Recorder
# No After=audioGuestBook.service: upstream's unit has After=multi-user.target,
# so depending on it creates an ordering cycle and systemd deletes our start
# job at boot. We don't need it anyway — the sidecar only polls GPIO and files.

[Service]
ExecStart=/usr/bin/python3 /home/admin/time-capsule-cam/video_recorder.py
WorkingDirectory=/home/admin/time-capsule-cam
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
# root: recordings dir and config.yaml are root-owned on the /opt-based
# images, and upstream's own services run as root too
User=root

[Install]
WantedBy=multi-user.target
